# deploy/nloop.service
[Unit]
Description=nloop — loop engineering engine + dashboard
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=exec
WorkingDirectory=/opt/nloop
# claude CLI ada di ~/.local/bin — PATH default systemd nggak mencakup itu
Environment=PATH=/root/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/opt/nloop/run.sh
Restart=on-failure
RestartSec=3
User=root

[Install]
WantedBy=multi-user.target
